# systemd/heleon-keepalive.service
[Unit]
Description=Hele-On Tracker keep-alive ping (Render free tier)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=Environment=HELEON_RENDER_URL=https://bus-tracker-a36o.onrender.com/healthz
ExecStart=/home/patrick.parks/code/scripts/keepalive.sh
StandardOutput=journal
StandardError=journal

# Hard timeout — Render cold start can take ~30 s, give it room.
TimeoutStartSec=40